# services/lighthouse-wiki_crawler.service
[Unit]
Description=Lighthouse Background Service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lighthouse/wiki_crawler.py
Restart=always

[Install]
WantedBy=multi-user.target
