# services/lighthouse-secret-folders.service
[Unit]
Description=Lighthouse Secret Folder Manager
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lighthouse/secret-folders/secret_folder_manager.py
Restart=always

[Install]
WantedBy=default.target
